# systemd/user/codex-master-watchdog.service
[Unit]
Description=Codex Master Fleet Watchdog
Documentation=man:systemd.service(5)

[Service]
Type=oneshot
Environment=CODEX_MASTER_MCP_INSTANCE_ID=codex-master-watchdog
CapabilityBoundingSet=
KeyringMode=private
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadWritePaths=%h/.local/state/codex-master-mcp %t
IPAddressDeny=any
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
UMask=0077
ExecStart=%h/.local/bin/codex-master-mcp watchdog all --idle-seconds 60 --poll-interval-seconds 15 --report-grace-seconds 15 --action interrupt --manage-unclaimed --quiet
